# deploy/ct-collect-tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cortex.collect-tick</string>

  <key>ProgramArguments</key>
  <array>
    <string>__VENV_PYTHON__</string>
    <string>-m</string>
    <string>cortex.collect_tick</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH_ENV__</string>
  </dict>

  <key>StartInterval</key>
  <integer>__COLLECT_INTERVAL_SEC__</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/ct-collect-tick.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/ct-collect-tick.log</string>
</dict>
</plist>
